# ops/systemd/knvox-admin-console.service.example
[Unit]
Description=KNVOX Admin Console
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/knvox-carrier
EnvironmentFile=/opt/knvox-carrier/.env
EnvironmentFile=-/etc/knvox/admin-console.env
ExecStart=/usr/bin/node /opt/knvox-carrier/apps/admin-console/server.js
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
